# ops/systemd/potion-health.service
[Unit]
Description=Potion Runner Health Snapshot
Wants=potion-runner.service

[Service]
Type=oneshot
User=dre
WorkingDirectory=/mnt/dre/potion-runner
ExecStart=/mnt/dre/potion-runner/venv/bin/python /mnt/dre/potion-runner/scripts/health_ping.py
