# ops/systemd/potion-runner-backup.service
[Unit]
Description=Potion Runner Backup Snapshot
Wants=potion-runner.service
After=network-online.target

[Service]
Type=oneshot
User=dre
WorkingDirectory=/mnt/dre/potion-runner
ExecStart=/mnt/dre/potion-runner/venv/bin/python /mnt/dre/potion-runner/scripts/backup_once.py
