# ops/systemd/potion-runner-backup.timer
[Unit]
Description=Potion Runner Backup Timer (02:30)

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
Unit=potion-runner-backup.service

[Install]
WantedBy=timers.target
