# ops/systemd/potion-runner-health.service
[Unit]
Description=Potion Runner Health Snapshot
Wants=potion-runner.service

[Service]
Type=oneshot
User=dre
WorkingDirectory=/opt/potion-runner
ExecStart=/opt/potion-runner/venv/bin/python /opt/potion-runner/scripts/health_ping.py
