# ops/systemd/potion-runner.service
[Unit]
Description=Potion Runner Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
User=dre
WorkingDirectory=/mnt/dre/potion-runner
Environment=PYTHONUNBUFFERED=1
ExecStart=/mnt/dre/potion-runner/venv/bin/python -m app.bot
Restart=always
RestartSec=5
WatchdogSec=60
StartLimitIntervalSec=300
StartLimitBurst=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
